# transmission-helper.service.example
[Unit]
Description=transmission-helper service

[Service]
ExecStart=/path/to/transmission-helper
Environment="TH_CONFIG_PATH=/path/to/config"

[Install]
WantedBy=multi-user.target
