# ansible/roles/systemd/templates/transmission-helper.service.j2
[Unit]
Description=transmission-helper service

[Service]
ExecStart={{ th_binary_path }}
Environment="TH_CONFIG_PATH={{ th_config_path }}"

[Install]
WantedBy=multi-user.target
