# ansible/roles/systemd/templates/transmission-helper.timer.j2
[Unit]
Description=Run transmission-helper regularly

[Timer]
Unit=transmission-helper.service
OnCalendar={{ timer_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
